# firefly-iii/ops/docker/cron/crontab
# Firefly III cron jobs - runs every minute to process scheduled tasks
# Artisan commands for Firefly III recurring tasks
* * * * * cd /var/www/html && php artisan schedule:run >> /var/log/cron/firefly.log 2>&1

# Generate recurring transactions - runs daily at 3 AM
0 3 * * * cd /var/www/html && php artisan firefly-iii:recurring-transactions >> /var/log/cron/recurring.log 2>&1

# Auto budget handling - runs daily at 2 AM  
0 2 * * * cd /var/www/html && php artisan firefly-iii:auto-budget >> /var/log/cron/budget.log 2>&1

# Email reports - runs weekly on Sunday at 8 AM
0 8 * * 0 cd /var/www/html && php artisan firefly-iii:send-reports >> /var/log/cron/reports.log 2>&1